# autostart/flameshot.desktop
[Desktop Entry]
Exec=flameshot
Name=flameshot
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
