# autostart/picom -b.desktop
[Desktop Entry]
Exec=picom -b
Name=picom -b
Type=Application
Version=1.0
