# autostart/pipewire-pulse.desktop
[Desktop Entry]
Exec=pipewire-pulse
Name=pipewire-pulse
Type=Application
Version=1.0
