# autostart/wireplumber.desktop
[Desktop Entry]
Exec=wireplumber
Name=wireplumber
Type=Application
Version=1.0
